# SeeTheWorld-Script.service
[Unit]
Description=Dump bing pictures

[Service]
ExecStart=/usr/bin/python3.9 -u /data/SeeTheWorld-Script/main.py